# google-drive/ai-sorter.service
[Unit]
Description=AI Drive Sorter - Inbox Auto-Router
After=network.target

[Service]
Type=oneshot
ExecStart=/home/tariqk/github/tariqk00/toolbox/google-drive/venv/bin/python3 /home/tariqk/github/tariqk00/toolbox/google-drive/drive_organizer.py --inbox --execute
WorkingDirectory=/home/tariqk/github/tariqk00/toolbox/google-drive
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
