# google-drive/ai-sorter.service
[Unit]
Description=AI Drive Sorter - Inbox Auto-Router
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/github/tariqk00/toolbox/google-drive/drive_organizer.py --inbox --execute
WorkingDirectory=%h/github/tariqk00/toolbox/google-drive
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
