# google-drive/drive-sorter.timer
[Unit]
Description=Run Google Drive Sorter Hourly

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
